# crontab
# Python Cron Jobs
# Format: minute hour day month weekday command
# All times in container timezone (set in docker-compose.yml)
#
# Scripts run in AUTO mode (no arguments) and automatically:
# - Find the last downloaded file
# - Download all missing files from last download date + 1 to yesterday
# - If no files exist, download last 7 days

# OTE-CR - Download day-ahead prices daily at 16:00
0 16 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_day_ahead_prices.py >> /var/log/cron.log 2>&1

# OTE-CR - Download imbalance prices every 2 hours
0 */2 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_imbalance_prices.py >> /var/log/cron.log 2>&1

# OTE-CR - Download intraday prices every 15 minutes
0 */15 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_intraday_prices.py >> /var/log/cron.log 2>&1

# OTE-CR - Download OTE daily settlement reports at 09:00
0 09 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 ote_production.py >> /var/log/cron.log 2>&1

# OTE-CR - Download OTE trade balance reports every 2 hours
0 */2 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 ote_trade_balance_downloader.py >> /var/log/cron.log 2>&1

# ENTSO-E - Unified imbalance data pipeline every 15 minutes
# Fetches imbalance prices (A85) and volumes (A86) for all areas
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_imbalance_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified load data pipeline every 15 minutes
# Fetches actual load (A65/A16) and day-ahead forecast (A65/A01) for all areas
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_load_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified generation data pipeline every 15 minutes
# Fetches actual generation per type (A75) for ALL areas: CZ, DE, AT, PL, SK
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_gen_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified cross-border flows pipeline every 15 minutes
# Fetches physical flows (A11) for all CZ border pairs (DE, AT, PL, SK)
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_flow_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified generation forecast pipeline every 15 minutes
# Fetches day-ahead wind/solar forecasts (A69) for all areas
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_forecast_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified balancing energy pipeline every 15 minutes
# Fetches activated aFRR/mFRR reserves (A84) for all areas
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_balancing_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified scheduled generation pipeline every 15 minutes
# Fetches day-ahead scheduled generation (A71) for all areas
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_scheduled_runner >> /var/log/cron.log 2>&1

# ENTSO-E - Unified scheduled cross-border flows pipeline every 15 minutes
# Fetches day-ahead scheduled exchanges (A09) for all CZ border pairs
*/15 * * * * export $(cat /etc/environment_for_cron | xargs) && cd /app && /usr/local/bin/python3 -m runners.entsoe_unified_sched_flow_runner >> /var/log/cron.log 2>&1

# ENTSO-E - XML cleanup daily at 02:00
# Retention policy: Delete XML files older than 1 day (keep only last 24 hours for verification)
0 02 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 entsoe/cleanup_old_xml.py >> /var/log/cron.log 2>&1

# Manual test command for OTE portal login:
# export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 ote_OTE-CR_login.py

# Empty line required at end of crontab file
