# crontab
# Python Cron Jobs
# Format: minute hour day month weekday command
# All times in container timezone (set in docker-compose.yml)
#
# Scripts run in AUTO mode (no arguments) and automatically:
# - Find the last downloaded file
# - Download all missing files from last download date + 1 to yesterday
# - If no files exist, download last 7 days

# TEST - Download day-ahead prices at 13:04
0 03 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_day_ahead_prices.py >> /var/log/cron.log 2>&1

# TEST - Download imbalance prices at 13:05
1 03 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_imbalance_prices.py >> /var/log/cron.log 2>&1

# TEST - Download intraday prices at 13:06
2 03 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 download_intraday_prices.py >> /var/log/cron.log 2>&1

# PRODUCTION - Download OTE daily settlement reports at 09:00
0 09 * * * export $(cat /etc/environment_for_cron | xargs) && cd /app/scripts && /usr/local/bin/python3 ote_production.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab file
